# docker/cron/crontab
#every minutes section
* * * * * php /app/bin/console parse:general
